# deploy/warehouse/jetson/warehouse-ros-bridge.service
[Unit]
Description=Warehouse ROS 2 Isaac Mapping Bridge
After=network-online.target
Wants=network-online.target

[Service]
Environment=WAREHOUSE_ROS_BRIDGE_HOST=0.0.0.0
Environment=WAREHOUSE_ROS_BRIDGE_PORT=8088
Environment=WAREHOUSE_ROS_CAPTURE_ROOT=/backend/storage/warehouse_ros
Environment=WAREHOUSE_ROS_PROFILE=isaac_ros_nvblox_stereo
Environment=WAREHOUSE_ROS_AUTOLAUNCH=1
Environment=WAREHOUSE_ROS_LAUNCH_CMD=ros2 launch warehouse_mapping_bridge isaac_warehouse_mapping.launch.py
ExecStart=/usr/bin/bash -lc 'source /opt/ros/humble/setup.bash && source /opt/warehouse_ros/install/setup.bash && ros2 run warehouse_mapping_bridge bridge_service'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

